# tools/namaz.service
[Unit]
Description=Namaz application
After=multi-user.target

[Service]
Type=simple
ExecStart=/home/mendel/namaz/baslat.sh > logs.xtx
WorkingDirectory=/home/mendel/namaz/
Restart=always

[Install]
WantedBy=multi-user.target

